# machines/rocky8/screenshot.timer
[Unit]
Description=Run screenshot service every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
